# apache/rootfs/etc/cron.d/nextcloud
*/5 * * * * www-data /beheer/scripts/nextcloud-cron.sh
