# mariadb/rootfs/etc/supervisor/conf.d/supervisor.conf
[supervisord]
user=root
nodaemon=true

[program:syslog-ng]
command=/etc/init.d/syslog-ng start

[program:mariadb]
command=service mariadb start
