# mail/rootfs/etc/supervisor/conf.d/supervisor.conf
[supervisord]
user=root
nodaemon=true

[unix_http_server]
file=/dev/shm/supervisor.sock
chmod=0700

[rpcinterface:supervisor]
supervisor.rpcinterface_factory = supervisor.rpcinterface:make_main_rpcinterface

[supervisorctl]
serverurl=unix:///dev/shm/supervisor.sock

[include]
files = /etc/supervisor/conf.d/*.conf

[program:syslog-ng]
command=/etc/init.d/syslog-ng start

[program:cron]
command=/usr/sbin/cron -f
autorestart=true

[program:postfix]
process_name=master
directory=/etc/postfix
command=/usr/sbin/postfix -c /etc/postfix start
startsecs=0
autorestart=false

[program:dovecot]
command=/usr/sbin/dovecot -c /etc/dovecot/dovecot.conf -F
autorestart=true

[program:opendkim]
command=/usr/sbin/opendkim -f -l -x /etc/opendkim/opendkim.conf
autorestart=true

[program:opendmarc]
command=/usr/sbin/opendmarc -f -l -c /etc/opendmarc/opendmarc.conf
autorestart=true

[program:spamassassin]
command=service spamd start
autorestart=true

[program:spamass-milter]
command=/etc/init.d/spamass-milter start
autorestart=true

#[program:fail2ban]
#command = service fail2ban start
